# polipo.conf
[program:polipo]
command=/start_polipo.sh
